# data/systemd/anaconda.target
[Unit]
Description=The anaconda installation program
